# system/brightnessd.service
[Unit]
After=

[Service]
ExecStart=/home/pi/carpc/system/set_brightness.sh

[Install]
WantedBy=default.target
